# etc/crontabs/root
*/5 * * * * /sbin/fan_ctrl.sh
*/1 * * * * sh /root/updateRouter.sh > /tmp/log/root 2>&1
